# api.service
[Unit]
Description=Dating Bot API Service
After=network.target postgresql.service

[Service]
Type=notify
User=your-user
Group=your-group
WorkingDirectory=/path/to/dating-bot-dev
Environment="PATH=/path/to/dating-bot-dev/venv/bin"

# Загрузка переменных окружения из .env файла
EnvironmentFile=/path/to/dating-bot-dev/.env

# Запуск через gunicorn
ExecStart=/path/to/dating-bot-dev/venv/bin/gunicorn -c gunicorn_config.py webapp:app

# Перезапуск при падении
Restart=always
RestartSec=10

# Ограничения
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
